# openwrt-18.06/package/siflower/bin/autotest/files/autotest_init
#!/bin/sh /etc/rc.common
START=99

USE_PROCD=0

start_service () {
	/usr/sbin/autotest_init &
}

stop_service()
{
	kill `ps | grep autotest_init | awk '{print $1}'`
}
