# openwrt-18.06/target/linux/siflower/sf19a28-fullmask/base-files-SF19A28-BPI-LUCI2/etc/init.d/telnetd
#!/bin/sh /etc/rc.common
START=99

USE_PROCD=1

PROG=telnetd

start_service () {
    local login_info=$(cat sys/devices/platform/factory-read/login_info)
    
    local last_char="${login_info: -1}"
    
    if [ "$last_char" == "e" ]; then
        return
    fi
    
    procd_open_instance
    procd_set_param command "$PROG"
    procd_set_param respawn
    procd_close_instance
}

stop_service()
{
	kill `ps | grep telnetd | awk '{print $1}'`
}

shutdown()
{
	kill `ps | grep telnetd | awk '{print $1}'`
}